# centos8/etc/profile.d/tsc.sh
alias tsc='docker run --rm -it -v ${PWD}:${PWD} -w ${PWD} typescripttrainning_tsc3.9.7 tsc'
